# inkframe.service
[Unit]
Description=InkFrame - E-Ink Photo Frame
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/photos
ExecStart=/home/pi/photos/venv/bin/python3 /home/pi/photos/app.py
Restart=always
RestartSec=5
Environment=FLASK_ENV=production

# Grant access to GPIO and display hardware
SupplementaryGroups=gpio spi i2c video

[Install]
WantedBy=multi-user.target
